# files/tuned-settings/etc/systemd/user/zodium-tuned-2-firmware-sync.path
[Unit]
Description=Watch tuned active profile for changes

[Path]
PathChanged=/etc/tuned/active_profile

[Install]
WantedBy=multi-user.target