# files/tuned-profiles/etc/systemd/user/zodium-firmware-2-tuned-sync.path
[Unit]
Description=Watch firmware platform_profile for changes

[Path]
PathChanged=/sys/firmware/acpi/platform_profile

[Install]
WantedBy=multi-user.target
