# files/zodium-settings/usr/lib/systemd/system/zodium-flatpak-repo.service
[Unit]
Description=Switch Flatpak remotes to Flathub
After=network-online.target
Wants=network-online.target
ConditionPathExists=!/var/lib/.flathub-migrated

[Service]
Type=oneshot
ExecStart=/usr/bin/bash -c '\
  flatpak remote-delete --system fedora 2>/dev/null || true; \
  flatpak remote-delete --system fedora-testing 2>/dev/null || true; \
  flatpak remote-add --system --if-not-exists flathub https://dl.flathub.org/repo/flathub.flatpakrepo; \
  flatpak remote-modify --system --enable flathub; \
  touch /var/lib/.flathub-migrated'
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
