# files/tuned-settings/usr/lib/udev/rules.d/09-zodium-platform-profile.rules
SUBSYSTEM=="platform-profile", ACTION=="change", \
RUN+="/usr/bin/systemctl start zodium-firmware-2-tuned-sync.service"