# files/tuned-settings/etc/systemd/system/zodium-tuned-2-firmware-sync.service
[Unit]
Description=Bidirectional sync between tuned and firmware

[Service]
Type=oneshot
ExecStart=/usr/libexec/zodium-tuned-sync.sh
