# files/tuned-settings/etc/systemd/system/zodium-firmware-2-tuned-sync.path
[Unit]
Description=Watch firmware platform_profile for changes

[Path]
PathModifiedc=/sys/firmware/acpi/platform_profile

[Install]
WantedBy=multi-user.target
