# files/base-config/usr/lib/udev/rules.d/09-tuned-platform-profile.rules
ACTION=="change", SUBSYSTEM=="platform_profile", RUN+="/usr/libexec/zodium-tuned-sync.sh"
